# systemd/polybot-dca.service.example
[Unit]
Description=Polymarket DCA Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/crypto-trading-bot
ExecStart=/root/crypto-trading-bot/bot-env/bin/python3 -m polybot.main_dca
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
Environment=DCA_BUDGET=100
Environment=DCA_INTERVAL_H=3
Environment=DCA_TOP_N=4
Environment=DCA_RESCAN_H=2
Environment=DCA_MIN_EDGE_PCT=0.50
Environment=DCA_NEG_STREAK_LIMIT=2
Environment=DCA_COIN_COOLDOWN_H=8
Environment=DCA_ROLLING_WINDOW=6
Environment=DCA_ROLLING_LOSS_LIMIT_EUR=-1.00
Environment=DCA_RISK_OFF_H=4
Environment=DCA_TAKE_PROFIT_PCT=0.03
Environment=DCA_STOP_LOSS_PCT=0.0
Environment=DCA_MAX_HOLD_SEC=1209600
Environment=DCA_MIN_NET_PROFIT_EUR=0.15
Environment=DCA_MAX_OPEN_POSITIONS=4
Environment=DCA_MAX_PAIR_EXPOSURE_EUR=25
Environment=DCA_MIN_CASH_RESERVE_EUR=5
Environment=DCA_TREND_FILTER_ENABLED=true
Environment=DCA_BTC_RISK_OFF_PCT=-2.0
Environment=DCA_ETH_RISK_OFF_PCT=-3.0
Environment=DCA_RECOVERY_TRIGGER_PCT=-5.0
Environment=DCA_RECOVERY_REVERSAL_PCT=0.8
Environment=DCA_RECOVERY_TICKET_EUR=7.0
Environment=DCA_RECOVERY_MAX_EXPOSURE_FACTOR=1.5

[Install]
WantedBy=multi-user.target
