# systemd/polybot-futures.service.example
[Unit]
Description=Polybot Futures (Der Hebler) — Paper Trading
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/crypto-trading-bot
ExecStart=/root/crypto-trading-bot/bot-env/bin/python3 -m polybot.main_futures
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1
Environment=FUT_BUDGET=100
Environment=FUT_SYMBOLS=PF_XBTUSD,PF_ETHUSD,PF_SOLUSD
Environment=FUT_INTERVAL_SEC=300
Environment=FUT_CANDLE_RESOLUTION=1h
Environment=FUT_FAST_EMA=9
Environment=FUT_SLOW_EMA=21
Environment=FUT_MOMENTUM_BARS=6
Environment=FUT_MIN_MOMENTUM_PCT=0.8
Environment=FUT_MIN_TREND_PCT=0.15
Environment=FUT_MIN_VOLUME_USD=10000000
Environment=FUT_POSITION_MARGIN_EUR=20
Environment=FUT_LEVERAGE=2
Environment=FUT_MAX_OPEN_POSITIONS=2
Environment=FUT_HARD_STOP_PCT=2
Environment=FUT_TAKE_PROFIT_PCT=5
Environment=FUT_TRAILING_ACTIVATION_PCT=2
Environment=FUT_TRAILING_DISTANCE_PCT=1
Environment=FUT_MAX_HOLD_H=24
Environment=FUT_COOLDOWN_H=4
Environment=FUT_MAX_SPREAD_PCT=0.10
Environment=FUT_MAX_ADVERSE_FUNDING_RATE=0.0002
Environment=FUT_TAKER_FEE_RATE=0.0005
Environment=FUT_MAINTENANCE_MARGIN_RATE=0.05
Environment=FUT_DAILY_LOSS_HALT_PCT=5
Environment=FUT_PAPER_MODE=true

[Install]
WantedBy=multi-user.target
