# systemd/polybot-battle-report.service.example
[Unit]
Description=Polybot Strategy Battle Report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/crypto-trading-bot
ExecStart=/root/crypto-trading-bot/bot-env/bin/python3 -m polybot.battle_report
Environment=PYTHONUNBUFFERED=1
