# systemd/polybot-pumpfun.service.example
[Unit]
Description=Polybot Pump.fun — Paper Trading Only
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/crypto-trading-bot
ExecStart=/root/crypto-trading-bot/bot-env/bin/python3 -m polybot.main_pumpfun
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1
Environment=PUMP_BUDGET=100
Environment=PUMP_POSITION_EUR=5
Environment=PUMP_MAX_OPEN_POSITIONS=2
Environment=PUMP_MIN_AGE_SEC=90
Environment=PUMP_MAX_AGE_SEC=21600
Environment=PUMP_MIN_MARKET_CAP_SOL=20
Environment=PUMP_MAX_MARKET_CAP_SOL=350
Environment=PUMP_MIN_CHANGE_PCT=10
Environment=PUMP_MAX_CHANGE_PCT=35
Environment=PUMP_MIN_TRADES=20
Environment=PUMP_MIN_BUY_SELL_RATIO=1.4
Environment=PUMP_STOP_LOSS_PCT=20
Environment=PUMP_TAKE_PROFIT_PCT=30
Environment=PUMP_TRAILING_STOP_PCT=15
Environment=PUMP_TRAIL_FLOOR_PCT=15
Environment=PUMP_MAX_HOLD_SEC=2700
Environment=PUMP_SLIPPAGE_PCT=2.5
Environment=PUMP_PLATFORM_FEE_PCT=1
Environment=PUMP_PAPER_MODE=true

[Install]
WantedBy=multi-user.target
