# systemd/polybot-meanrev.service.example
[Unit]
Description=Polybot Mean-Reversion (Der Contrarian) — Paper Trading
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/crypto-trading-bot
ExecStart=/root/crypto-trading-bot/bot-env/bin/python3 -m polybot.main_meanrev
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1
Environment=REV_BUDGET=100
Environment=REV_INTERVAL_H=1
Environment=REV_ENTRY_DROP_PCT=8.0
Environment=REV_RSI_PERIOD=14
Environment=REV_RSI_MAX=30
Environment=REV_BOLLINGER_ENABLED=true
Environment=REV_BOLLINGER_PERIOD=20
Environment=REV_BOLLINGER_STDDEV=2.0
Environment=REV_STOCHASTIC_ENABLED=true
Environment=REV_STOCHASTIC_PERIOD=14
Environment=REV_STOCHASTIC_MAX=20
Environment=REV_CONFIRM_PCT=0.5
Environment=REV_POSITION_EUR=15
Environment=REV_MAX_OPEN_POSITIONS=3
Environment=REV_TAKE_PROFIT_PCT=4.0
Environment=REV_STOP_LOSS_PCT=5.0
Environment=REV_MAX_HOLD_H=96
Environment=REV_COOLDOWN_H=12
Environment=REV_PAPER_MODE=true

[Install]
WantedBy=multi-user.target
