# systemd/polybot-futures.service.example
[Unit]
Description=Polybot Futures Grid (Der Treppensteiger Turbo) — Paper Trading
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/crypto-trading-bot
ExecStart=/root/crypto-trading-bot/bot-env/bin/python3 -m polybot.main_futures
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1
Environment=FUT_BUDGET=1000
Environment=FUT_LEVERAGE=2
Environment=FUT_ORDER_MARGIN_EUR=15
Environment=FUT_GRID_STEP_PCT=0.8
Environment=FUT_TAKE_PROFIT_PCT=1.1
Environment=FUT_MAX_SAFETY_ORDERS=50
Environment=FUT_MAINTENANCE_MARGIN_PCT=5
Environment=FUT_MARGIN_GUARD_RATIO=1.25
Environment=FUT_TAKER_FEE_RATE=0.0005
Environment=FUT_FUNDING_RATE_8H=0.0001
Environment=FUT_SCAN_INTERVAL_SEC=30
Environment=FUT_PAPER_MODE=true

[Install]
WantedBy=multi-user.target
